# infra/systemd/haven-dashboard.service
[Unit]
Description=Haven dashboard (SvelteKit SSR via adapter-node)
Documentation=https://github.com/michael/haven
After=network-online.target haven-backend.service
Wants=network-online.target
Requires=haven-backend.service

[Service]
Type=simple
User=haven
Group=haven
WorkingDirectory=/opt/haven/apps/dashboard
EnvironmentFile=/etc/haven/.env
# 0.0.0.0 so a reverse proxy on a different bridge (e.g. Docker Caddy)
# can reach the SSR. Override to 127.0.0.1 if running everything on the
# host network and you want to limit exposure.
Environment=HOST=0.0.0.0
Environment=PORT=3000
# ORIGIN is sourced from /etc/haven/.env (set it to the public URL —
# SvelteKit needs it for CSRF on POST/form actions).
ExecStart=/usr/local/bin/bun run build/index.js
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true

[Install]
WantedBy=multi-user.target
